# scripts/systemd/nimsforest.service
[Unit]
Description=NimsForest Event Orchestration System
After=network-online.target mnt-forest.mount
Wants=network-online.target
Documentation=https://github.com/yourusername/nimsforest

[Service]
Type=simple
User=forest
Group=forest
WorkingDirectory=/var/lib/nimsforest

# Environment configuration
# Environment="JETSTREAM_DIR=/var/lib/nimsforest/jetstream"
# Environment="DEMO=false"

# Service execution
ExecStart=/opt/nimsforest/bin/forest
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=always
RestartSec=5
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nimsforest

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/nimsforest /var/log/nimsforest /mnt/forest
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Resource limits
LimitNOFILE=65536
TasksMax=4096

[Install]
WantedBy=multi-user.target
